# playbooks/roles/riak/files/etc_init_riak_explorer.conf
# riak_explorer - Riak Explorer

description     "Riak Explorer"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec /usr/local/riak_explorer/bin/riak_explorer start
